# recipes-core/systemd/systemd-firmware/firmware.service
[Unit]
Description=Firmware Log Configuration for Realtek Avengers SoC

[Service]
ExecStart=/usr/bin/fwdbg 0

[Install]
WantedBy=graphical.target
