# recipes-graphics/wayland/weston-init/52-drm.rules
ACTION=="change", KERNEL=="card[0-9]", SUBSYSTEM=="drm", RUN+="/usr/bin/start-weston.sh"
